# contrib/disk_monitor/etc/cron.hourly/disk_monitor
#!/bin/sh

/root/disk_monitor.sh

exit 0

